# config/linux_conky/.config/autostart/conky.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/bash -l -c "sleep 5 && conky"
Name=Conky
Comment=Start Conky at login
